# led5off.timer
[Unit]
Description=Turn off status LEDs 90 seconds after boot

[Timer]
OnBootSec=90sec
Unit=led5off.service

[Install]
WantedBy=timers.target
